# os/profiles/axiom-desktop-files/user_scripts/update_axiom/update_checker/service/update_checker.service
[Unit]
Description=Axiom Update Checker Service

[Service]
Type=oneshot
ExecStart=%h/user_scripts/update_axiom/update_checker/update_checker.sh --num
