# os/profiles/axiom-desktop-files/user_scripts/axiom_system/control_center/service/axiom.service
[Unit]
Description=Axiom Control Center Daemon
# Wait for Hyprland to be ready to avoid "Display not found" errors
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=com.github.axiom.controlcenter

# CRITICAL: Forces Python to flush logs immediately for real-time debugging
Environment=PYTHONUNBUFFERED=1

# ENABLE RAM TRACKING
MemoryAccounting=yes

# The %h specifier automatically inserts your home directory path
ExecStart=%h/user_scripts/axiom_system/control_center/axiom_control_center.py --gapplication-service

# Restart automatically if it crashes
Restart=always
RestartSec=3

# Logs will show up as "axiom[PID]: ..."
SyslogIdentifier=axiom

# Better process scheduling group
Slice=app.slice

[Install]
WantedBy=graphical-session.target
